# s6/services.d/sks-stats/run
#!/bin/execlineb -P
foreground
{
  s6-svwait /var/run/s6/services/sks-db
}
loopwhilex foreground {
    s6-svc -2 /var/run/s6/services/sks-db
} sleep 1800
