# s6/services.d/sks-db/run
#!/bin/execlineb -P
/usr/local/bin/big-sks -stdoutlog -basedir /var/lib/sks db
